# examples/oracle_pull/templates/launchd_pull_schedule.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-council.oracle-pull-preview.disabled-example</string>
  <key>Disabled</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/echo</string>
    <string>Review scripts/run_oracle_pull_smoke.sh and configure placeholders before manual use. 실제 주문 없음. order_execution_allowed=false.</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
</dict>
</plist>
